# external/sketch/ansible/roles/disable-ipv6/files/restore-ip6tables.service
[Unit]
Description=Restore IPv6 iptables rules

[Service]
Type=oneshot
ExecStart=/sbin/ip6tables-restore /etc/ip6tables/rules.v6
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
